# deploy/ansible/roles/deploy/templates/node_server.conf
description "node server"

start on startup
stop on shutdown
respawn

script
  /usr/bin/node /mnt/site/server/server.js
end script
